# searxng/searxng.container
[Unit]
Description=SearXNG Instance
After=network-online.target

[Container]
Image=docker.io/searxng/searxng:2026.5.8-d8ab61a9e
ContainerName=searxng-core
Environment=SEARXNG_HOST=127.0.0.1
Environment=SEARXNG_PORT=9000
PublishPort=127.0.0.1:9000:9000
Volume=/home/alex/git/agent-sandbox/searxng/core-config:/etc/searxng:Z
Volume=searxng_core-data:/var/cache/searxng

HealthCmd=curl -sf http://localhost:9000/healthz
HealthInterval=30s
HealthTimeout=10s
HealthStartPeriod=60s

[Service]
Restart=on-failure
RestartSec=5
TimeoutStartSec=120
TimeoutStopSec=90

[Install]
WantedBy=default.target
